# tapd-backup.service
[Unit]
Description=Taproot Assets Backup Monitor (%NETWORK%)
After=network.target tapd.service
Wants=network.target

[Service]
Type=simple
ExecStart=%INSTALL_DIR%/tapd-backup-wrapper
Restart=on-failure
RestartSec=30

# User/Group for system service
%USER_GROUP_SECTION%

# Environment
EnvironmentFile=%CONFIG_DIR%/config
Environment="NETWORK=%NETWORK%"
Environment="CONFIG_DIR=%CONFIG_DIR%"

# Credential loading
%CREDENTIAL_SECTION%

# Security hardening
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=yes
ReadWritePaths=%CONFIG_DIR% /tmp
ReadOnlyPaths=%TAPD_DATA_DIR%

# Group access for tapd files
SupplementaryGroups=lndbackup

# Resource limits
MemoryMax=256M
CPUQuota=20%

[Install]
WantedBy=%WANTED_BY%